# src/templates/ai.openclaw.dashboard.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.dashboard</string>

    <key>ProgramArguments</key>
    <array>
        <string>__GOOSE_NODE_PATH__</string>
        <string>__GOOSE_HOME__/dashboard/server.js</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>WorkingDirectory</key>
    <string>__GOOSE_HOME__/dashboard</string>

    <key>StandardOutPath</key>
    <string>__GOOSE_HOME__/logs/dashboard-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__GOOSE_HOME__/logs/dashboard-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__GOOSE_USER_HOME__</string>
        <key>PATH</key>
        <string>__GOOSE_PATH__</string>
        <key>NODE_PATH</key>
        <string>__GOOSE_NPM_ROOT__</string>
        <key>NODE_ENV</key>
        <string>production</string>
        <key>PORT</key>
        <string>3333</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
