# src/templates/ai.openclaw.learn-daily.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.openclaw.learn-daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>__GOOSE_NODE_PATH__</string>
        <string>__GOOSE_HOME__/pipeline/learn-daily.js</string>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>WorkingDirectory</key>
    <string>__GOOSE_HOME__</string>

    <key>StandardOutPath</key>
    <string>__GOOSE_HOME__/logs/learn-daily-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__GOOSE_HOME__/logs/learn-daily-stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>__GOOSE_USER_HOME__</string>
        <key>PATH</key>
        <string>__GOOSE_PATH__</string>
        <key>NODE_PATH</key>
        <string>__GOOSE_NPM_ROOT__</string>
    </dict>
</dict>
</plist>
